# scripts/com.sentra.kolrefresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sentra.kolrefresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/luthfiamri/sentra-portal/scripts/refresh_kol.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/luthfiamri/sentra-portal</string>

    <!-- Monday (Weekday=1) at 12:00. If Mac is asleep, launchd fires on next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/luthfiamri/sentra-portal/scripts/refresh_kol.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/luthfiamri/sentra-portal/scripts/refresh_kol.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
